# com.user.imessage_relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.imessage_relay</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/nishiokamahiro/Desktop/Antigravity/iMessage soujushin/imessage_relay.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/nishiokamahiro/Library/Logs/imessage_relay.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/nishiokamahiro/Library/Logs/imessage_relay.log</string>
</dict>
</plist>
